# sfs-jam-monitor.service
[Unit]
Description=BTT SFS v2.0 Jam Monitor (Prusa Core One)
After=network-online.target
Wants=network-online.target

After=dev-ttyACM0.device dev-ttyUSB0.device
Wants=dev-ttyACM0.device dev-ttyUSB0.device

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sfs-jam-monitor

# IMPORTANT: Prefer /dev/serial/by-id/... when possible for stable device naming.
ExecStart=/usr/bin/python3 /home/pi/sfs-jam-monitor/sfs-monitor.py -p /dev/ttyACM0 --motion-gpio 26 --auto-reset --reset-pulses 1.5 --quiet-temps --runout-gpio 27

Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=true
SystemCallArchitectures=native

DevicePolicy=closed
DeviceAllow=/dev/ttyACM0 rw
DeviceAllow=/dev/ttyUSB0 rw
DeviceAllow=/dev/gpiomem rw

ReadWritePaths=/var/log /var/lib/node_exporter

[Install]
WantedBy=multi-user.target
